# vendor/gitlab.com/mergetb/tech/nex/systemd/nex-dhcpd.service
[Unit]
Description=nex-dhcpd
Documentation=https://gitlab.com/mergetb/nex

[Service]
ExecStart=/usr/bin/nex-dhcpd

[Install]
WantedBy=multi-user.target

