# vendor/gitlab.com/mergetb/tech/nex/systemd/nexd.service
[Unit]
Description=nexd
Documentation=https://gitlab.com/mergetb/nex

[Service]
ExecStart=/usr/bin/nexd

[Install]
WantedBy=multi-user.target

