# vendor/gitlab.com/mergetb/nex/coredns.service
[Unit]
Description=coredns
Documentation=https://coredns.io

[Service]
ExecStart=/usr/bin/coredns \
  -conf /etc/coredns/Corefile \
  -dns.port 53

[Install]
WantedBy=multi-user.target

